# config.dist/gunicorn-sixlinks.conf
description "Gunicorn for Sixlinks"

start on runlevel [2345]
stop on runlevel [!2345]
#Send KILL after 20 seconds
kill timeout 20

script
exec /root/.virtualenvs/sixlinks/bin/python /var/www/sixlinks.git/sixlinks_web/manage.py run_gunicorn 0.0.0.0:8001 -w 3 --preload 
end script

respawn
